# .config/autostart/arcade.desktop
[Desktop Entry]
Type=Application
Name=Arcade
Exec=/home/arcade/arcade-install.sh